# docker/production/supervisord.conf
; ===================================================================
;   Atlas - Supervisord (PID 1 do container)
;   Orquestra nginx + php-fpm. Se um cair, derruba tudo (Railway reinicia).
; ===================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
; Aguarda 1s pra php-fpm subir antes de iniciar nginx
command=sh -c "sleep 1 && nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Se qualquer um sair com codigo != 0 ou exit normal, encerra tudo
; Railway re-inicia o container automaticamente
[eventlistener:processes]
command=sh -c "printf 'READY\n' && while read -r line; do kill -SIGQUIT \$PPID; done"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
